# Library/io.git.notes_add.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- File Location: $HOME/Library/LaunchAgents/io.git.notes_add.plist-->
<plist version="1.0">
    <dict>
        <key>Label</key>
        <string>io.git.notes_add</string>
        <key>ProgramArguments</key>
        <array>
            <string>git</string>
            <string>add</string>
            <string>.</string>
        </array>
        <key>WorkingDirectory</key>
        <string>/Users/bhavesh/NotesLogSeq</string>
        <key>EnvironmentVariables</key>
        <dict>
            <key>GIT_CURL_VERBOSE</key>
            <string>1</string>
            <key>GIT_TRACE_PACKET</key>
            <string>1</string>
            <key>GIT_TRACE</key>
            <string>1</string>
            <key>PATH</key>
            <string>/Users/bhavesh/.nix-profile/bin:/nix/var/nix/profiles/default/bin</string>
        </dict>
        <key>StartInterval</key>
        <integer>3480</integer>
        <key>StandardOutPath</key>
        <string>/Users/bhavesh/Library/Logs/git.notes/launchd-stdout.log</string>
        <key>StandardErrorPath</key>
        <string>/Users/bhavesh/Library/Logs/git.notes/launchd-stderr.log</string>
    </dict>
</plist>